# vast-vite.service
[Unit]
Description=VAST Vite Dev Server (port 8080)
After=network.target

[Service]
Type=simple
User=claude
WorkingDirectory=/home/claude/vast
ExecStart=/usr/bin/npx vite --host 0.0.0.0 --port 8080 --strictPort
Restart=always
RestartSec=2
Environment=NODE_ENV=development
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
