# vast-api.service
[Unit]
Description=VAST Screenshot API (port 8081)
After=network.target

[Service]
Type=simple
User=claude
WorkingDirectory=/home/claude/vast
Environment=PORT=8081
ExecStart=/usr/bin/node /home/claude/vast/api-server.cjs
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
